# autorun.conf
[program:autorun]
autorstart=true
directory = /home/pi/share/gitShare/supervisord-conf/autorun/
command = python app.py
