# raspberrypi.conf
[program:raspberrypi]
command=python /mnt/400gb/pylove/raspberrypi.py
autorstart=true 
stdout_logfile=/tmp/raspberrypi
